# garage_supervisor.conf
[program:garage]
command = /usr/bin/env uwsgi /home/pi/projects/uwsgi.ini
